# ludoc-os/systemd/ludoc-context-server.service
[Service]
Type=simple
ExecStart=/usr/bin/bun /opt/ludoc/src/index.ts
MemoryLimit=4G
MemoryMax=4G
CPUQuota=200%
